# data/systemd/experienced.service
[Unit]
Description=GNOME X Experience Daemon
Documentation=https://github.com/gnomex/gnome-x
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/experienced
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
